# deploy/supervisor/booking-platform-worker.conf
[program:booking-platform-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/booking-platform/artisan queue:work redis --queue=default,webhooks,notifications --sleep=1 --tries=5 --backoff=5 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/booking-platform/storage/logs/worker.log
stopwaitsecs=3600
